# project/packages/debian/DEBIAN/postinst
#!/bin/sh
update-rc.d mswitch multiuser
